# .crontab
* * * * * php artisan  drandru:add-file
